# build/linux-x64/beta-censor-server.service
[Unit]
Description=Beta Censoring Server

[Service]
Type=notify
ExecStart=/usr/local/bin/censor-core-server/BetaCensor.Server

[Install]
WantedBy=multi-user.target